# .config/autostart/xdg-desktop-portal-gtk.desktop
[Desktop Entry]
Hidden=true
Name=desktop-portal-gtk
Exec=/usr/lib/xdg-desktop-portal-gtk ; /usr/lib/xdg-desktop-portal --replace
